# etc/systemd/user/setup-esteid-pkcs11.service
[Service]
ExecStart=/opt/bin/setup-user-esteid-pkcs11.sh
Type=oneshot
RemainAfterExit=true
